# flask/flask-boost/test/deploy/supervisor.conf
[program:test]
command=/var/www/test/venv/bin/gunicorn -c deploy/gunicorn.conf wsgi:app
directory=/var/www/test
user=root
autostart=true
autorestart=true
environment = MODE="PRODUCTION"